# debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        # Enable and start the service
        if [ -d /run/systemd/system ]; then
            systemctl daemon-reload
            systemctl enable halos-mdns-publisher.service
            systemctl start halos-mdns-publisher.service || true
        fi
        ;;
esac

#DEBHELPER#

exit 0
